# scripts/systemd/grok-credential-checker.timer
[Unit]
Description=Run Grok credential checker every 5 minutes
Documentation=file:///opt/CLIProxyAPI/scripts/README-grok-credential-checker.md

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=grok-credential-checker-oneshot.service

[Install]
WantedBy=timers.target
