# scripts/systemd/grok-credential-checker.service
[Unit]
Description=CLIProxyAPI Grok free credential checker (daemon)
Documentation=file:///opt/CLIProxyAPI/scripts/README-grok-credential-checker.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cliproxy
Group=cliproxy
WorkingDirectory=/opt/CLIProxyAPI
EnvironmentFile=-/etc/cliproxy/grok-credential-checker.env
ExecStart=/usr/bin/python3 /opt/CLIProxyAPI/scripts/grok_credential_checker.py \
  --daemon \
  --apply \
  --interval-seconds 300 \
  --target-active 500 \
  --quota-limit-tokens 1000000 \
  --reset-window-hours 24 \
  --concurrency 32 \
  --batch-size 100 \
  --max-probes-per-cycle 50 \
  --state-file /var/lib/cliproxy/grok_credential_checker_state.json \
  --log-level INFO
Restart=on-failure
RestartSec=30
# Hardening (adjust if paths differ)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/cliproxy

[Install]
WantedBy=multi-user.target
